# minim-openwifi-scripts/files/etc/uci-defaults/34_reject_wan_input
# by default, reject incoming connections to wan
uci set firewall.@zone[1].input=REJECT
uci commit firewall
