# minim-openwifi-scripts/files/etc/uci-defaults/35_disable_maverick
# disable the maverick service
/etc/init.d/maverick disable
exit 0
