# minim-openwifi-scripts/files/etc/uci-defaults/36_disable_firstcontact_config
# If ucentral has configured this device, then exit immediately
active=$(readlink  /etc/ucentral/ucentral.active)
[ -n "$active" -a "$active" != "/etc/ucentral/ucentral.cfg.0000000001" ] && return 0

# We don't want /etc/init.d/firstcontact to apply the openwifi out-of-box config, so
# empty the config so it will harmlessly fail to apply
[ -f /etc/ucentral/ucentral.cfg.0000000001 ] && echo "{}" > /etc/ucentral/ucentral.cfg.0000000001

# We don't want to blow away the default wireless interfaces we created in
# 31_default_wireless, so skip 99-ucentral-wifi (which does exactly that)
rm -f /etc/uci-defaults/99-ucentral-wifi

# We also don't want to blow away the default openwrt network config,
# so skip 99-ucentral-network (which does exactly that)
rm -f /etc/uci-defaults/99-ucentral-network

# We also don't want to blow away the default openwrt uhttpd config,
# so skip 99-ucentral-uhttpd (which does exactly that)
rm -f /etc/uci-defaults/99-ucentral-uhttpd

# No srsly, don't delete the network config
sed -i '/rm \/etc\/config\/network/d' /etc/uci-defaults/zzz-ucentral

# Don't overwrite the uhttpd, firewall, dhcp, and dropbear configs
sed -i '/cp \/rom/d' /etc/uci-defaults/zzz-ucentral

exit 1
