# minim-openwifi-scripts/files/etc/uci-defaults/37_enable_flow_offloading
# enable software flow offloading
uci set firewall.@defaults[0].flow_offloading=1
uci commit firewall
