# minim-openwifi-scripts/files/etc/uci-defaults/38_enable_usb_automount_config
# setting USB auto mount

uci set fstab.@global[0].anon_mount=1


