# minim-openwifi-scripts/files/etc/uci-defaults/39_enable_upnp_enable_config
# setting to upnp enable configuration

uci set upnpd.config.enabled=1


