# minim-openwifi-scripts/files/etc/init.d/initialize_leds
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2021 OpenWrt.org

START=50

boot() {
    /sbin/leds.sh led_boot
}

